# build/pingmon.service
[Unit]
Description=Network connection monitoring
After=network-online.target

[Service]
ExecStart=/usr/local/bin/pingmon
ConfigDirectory=/usr/local/etc/
WorkingDirectory=/var/local/pingmon
Restart=on-success
RestartSec=10

[Install]
WantedBy=default.target
